# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=appuser

; ── FastAPI backend ───────────────────────────────────────────────────────────
[program:backend]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/app/logs/backend.log
stderr_logfile=/app/logs/backend.log
stdout_logfile_maxbytes=10MB
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ── Streamlit frontend ────────────────────────────────────────────────────────
; Waits 5 seconds for the backend to be ready before starting
[program:frontend]
command=streamlit run frontend/streamlit_app.py
    --server.port=7860
    --server.address=0.0.0.0
    --server.headless=true
    --server.enableCORS=false
    --server.enableXsrfProtection=false
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/app/logs/frontend.log
stderr_logfile=/app/logs/frontend.log
stdout_logfile_maxbytes=10MB
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",BACKEND_URL="http://localhost:8000"

; ── Process group ─────────────────────────────────────────────────────────────
[group:app]
programs=backend,frontend
